# systemd/check_chromium_versions.service
[Unit]
Description=Run check_chromium_versions

[Service]
Type=oneshot
ExecStart=/var/local/chromium-ci-trigger/check_chromium_versions.sh
